# ui/elata-news/elata-nextjs.service
[Unit]
Description=Elata Next.js Application
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/elata-vsm-system-4/ui/elata-news
Environment=NODE_ENV=production
Environment=PORT=3000

ExecStartPre=/usr/bin/npm install
ExecStartPre=/usr/bin/npm run build
ExecStart=/usr/bin/npm start

Restart=always
RestartSec=10
StandardOutput=append:/var/log/elata-nextjs.log
StandardError=append:/var/log/elata-nextjs.error.log

[Install]
WantedBy=multi-user.target